# batcountry.target
[Unit]
Description=batcountry target
StopWhenUnneeded=yes
